# root/etc/services.d/tftp/run
#!/usr/bin/with-contenv bash

/usr/sbin/in.tftpd \
	--foreground --listen --user abc --secure "${PORT_RANGE:+--port-range $PORT_RANGE}" /config/menus
