# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
exec /usr/sbin/nginx -c /config/nginx/nginx.conf